# etc/systemd/hegg-mqtt.service
[Unit]
Description=Hegg energy monitor — Home Assistant MQTT exporter
Documentation=https://github.com/yourorg/hegg
After=network-online.target hegg-collector.service
Wants=network-online.target

[Service]
Type=simple
User=hegg
Group=hegg
WorkingDirectory=/opt/hegg
Environment=HEGG_DB=/var/lib/hegg/hegg.db
ExecStart=/opt/hegg/.venv/bin/python /opt/hegg/hegg_mqtt.py \
    --db ${HEGG_DB}
EnvironmentFile=-/etc/hegg/hegg.conf
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
